# deploy/launchd/com.transfr.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent: the Cloudflare tunnel exposing the local API to the iOS app.

  TEMPLATE -- install.sh fills __CLOUDFLARED__ and __HOME__ and installs the
  rendered copy to ~/Library/LaunchAgents/.

  Runs the NAMED tunnel `transfr`, which serves a STABLE hostname
  (api.trans-fr.com) that survives restarts/reboots. It reads its routing and
  credentials from ~/.cloudflared/config.yml + the tunnel's <uuid>.json (both
  outside the repo, created once via `cloudflared tunnel login/create/route dns`).

  Fallback with no domain -- a QUICK tunnel (random URL, changes each restart):
      <string>__CLOUDFLARED__</string><string>tunnel</string>
      <string>--url</string><string>http://localhost:5001</string>
  then read the URL from tunnel.err.log.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.transfr.tunnel</string>

    <key>ProgramArguments</key>
    <array>
        <string>__CLOUDFLARED__</string>
        <string>tunnel</string>
        <string>run</string>
        <string>transfr</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/transfr/tunnel.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/transfr/tunnel.err.log</string>
</dict>
</plist>
